# misc/docker-demo/ableC-cpp.supervisord.conf
[supervisord]
nodaemon = true

[program:ableC]
command = java -jar ./ableC-monto.jar 10000

[program:cpp]
command = ./monto3-cpp
